# openswarm.service
[Unit]
Description=OpenSwarm Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/openswarm
EnvironmentFile=/root/openswarm/.env
ExecStart=/root/openswarm/.venv/bin/python discord_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
